# scripts/systemd/weppcloud-watchdog.service
[Unit]
Description=WEPPcloud Memory Watchdog

[Service]
Type=oneshot
ExecStart=/workdir/wepppy/scripts/memory-watchdog.sh
User=root
